# mqtt.crontab
1 * * * * /home/erwin2019/mqtt.sh
55 * * * * pkill -15 sh
55 * * * * pkill -15 mosquitto_sub
